# server/garmin-sync.timer
[Unit]
Description=Daily Garmin Sync Timer (3 AM local time)
Requires=garmin-sync.service

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
AccuracySec=5m

[Install]
WantedBy=timers.target
